# udiskie/debian/postinst
#!/bin/bash

# Allow users to access /media directory
if [ ! -f /etc/skel/media ];
then
	ln -s /media /etc/skel/
fi

update-rc.d udiskie defaults
service udiskie restart
