# yunohost-cli/debian/postinst
#!/bin/bash

update-rc.d yunohost-api defaults
