# moulinette/debian/postinst
#!/bin/bash

mkdir -p /usr/lib/moulinette
mkdir -p /usr/share/moulinette/actionsmap

update-rc.d yunohost-api defaults
service yunohost-api restart
