# glances/debian/postrm
#!/bin/bash

service glances stop
update-rc.d glances remove
